# deploy/systemd/trendradar-mcp.service
[Unit]
Description=TrendRadar MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trendradar
Group=trendradar
WorkingDirectory=/opt/trendradar
Environment=HOME=/var/lib/trendradar
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/uv run python -m mcp_server.server --transport http --host 127.0.0.1 --port 3333
Restart=always
RestartSec=5
KillSignal=SIGINT
SuccessExitStatus=130
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
