# deploy/systemd/trendradar-wsj-backfill.service
[Unit]
Description=One-shot current-list WSJ backfill to Feishu Docx
After=network-online.target bpc-server.service wsj-cn-rss-bridge.service
Wants=network-online.target
Requires=bpc-server.service wsj-cn-rss-bridge.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/TrendRadar
Environment=HOME=/home/ubuntu
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/trendradar/wsj-delivery.env
StateDirectory=trendradar
UMask=0077
ExecStart=/home/ubuntu/.local/bin/uv run python -m trendradar --wsj-delivery --backfill-current --drain
SuccessExitStatus=75
TimeoutStartSec=90min
